# soc/container-files-soc/etc/supervisor.d/snmptrapd.conf
[program:snmptrapd]
command = /usr/sbin/snmptrapd -f -Lo -n -M+/etc/soc/snmp/mibs -m+ALL
#command = /usr/sbin/snmptrapd -Lsd -f
autorestart = true
