# soc/container-files-soc/etc/supervisor.d/soc-agent.conf
[program:soc-agent]
command = /config/ds.sh /tmp/zabbix_agentd.pid sudo -u soc /usr/local/sbin/zabbix_agentd -f -c /usr/local/etc/soc_agentd.conf
autorestart = true
