# soc/container-files-soc/etc/supervisor.d/soc-server.conf
[program:soc-server]
command = /config/ds.sh /var/run/zabbix_server.pid zabbix_server -c /usr/local/etc/soc_server.conf
autorestart = true
