# soc/container-files-soc/etc/supervisor.d/soc-agent.conf
[program:soc-agent]
command = /config/ds.sh /tmp/soc_agentd.pid sudo -u soc /usr/local/sbin/soc_agentd -f -c /usr/local/etc/soc_agentd.conf
autorestart = true
