# soc/container-files-soc/etc/supervisor.d/soc-server.conf
[program:soc-server]
command = /config/ds.sh /tmp/soc_server.pid soc_server -c /usr/local/etc/soc_server.conf
autorestart = true
